# scripts/email-ip-notifier/send-ip-email.service
[Unit]
Description=Send IP Address via Email on Boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 INSTALL_DIR/send_ip_email.py
User=INSTALL_USER
WorkingDirectory=INSTALL_DIR
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
